# content/etcupstart/init/progress-services.conf
#!upstart
description "progress-filesystems-mounted"

start on filesystem 

task

script

	[ -e /run/nosplash ] && exit 0

	splash --msgtxt="loading..." || true 

end script
