# content/etc/init/failsafe-boot-started.conf
start on failsafe-boot

pre-start exec start -n networking
