# content/etc/init/failsafe-boot-started.conf
start on failsafe-boot

pre-start script
    start -qn networking || [ -d /run/network/static-network-up-emitted ] || ifup --force eth0 || :
end script 
