# content/etc/init/hwclock-systemd-disable.conf
start on started hwclock.sh

task

script
    mkdir -p /run/udev
    touch /run/udev/hwclock-set
end script
