# content/etc/init/iwconfig.conf
start on net-device-added INTERFACE=wlan0 or net-device-added INTERFACE=wlan1

instance $INTERFACE

exec iwconfig $INTERFACE power off
