# content/etc/init/motd.conf
start on runlevel [2]


task

exec /etc/cron.daily/updates-available || :

