# content/etc/init/ntp.conf
#!upstart

description "waits for ntpd to be started, then emits ntp event"

start on started networking 

emits ntp

pre-start script
    ntpd -g -q
    initctl emit -n ntp

    ! pgrep ntpd > /dev/null && service ntp start || true &
end script

post-stop script
    service ntp stop
end script
