# content/etc/init/openbsd-inetd-reload.conf
start on net-device-up and started openbsd-inetd

task

pre-start script
	sleep 1
end script

exec /sbin/reload openbsd-inetd
