# content/etc/init/openbsd-inetd-reconf-on-netchange.conf
start on net-device-up

task

exec /sbin/reload openbsd-inetd
