# content/etc/init/progress-mounting.conf
#!upstart
description "progress on events"

start on mounting or net-device-up or static-network-up
stop on mounted MOUNTPOINT="$MOUNTPOINT" TYPE=$TYPE DEVICE=$DEVICE

emits mounted

env IFACE=
env DEVICE
env INTERFACE=
env MOUNTPOINT=
env TYPE=
env SNU=y

instance $SNU$MOUNTPOINT$IFACE$TYPE

normal exit TERM

nice -2

pre-start script
    set +e
    case $UPSTART_EVENTS in
        mounting)
            ! mountpoint -q "$MOUNTPOINT" || { stop; exit 0; }
            echo "$(cat /proc/uptime) $UPSTART_EVENTS $MOUNTPOINT $TYPE $DEVICE" >> /run/upstart-mnt.log
            [ -n "$MOUNTPOINT" ] || stop
            ;;
        net-device-up|static-network-up)
            [ -n "$IFACE" ] || IFACE=$SNU
            echo "$(cat /proc/uptime) interface $IFACE is up" >> /run/upstart-mnt.log
            stop
            ;;
    esac
end script

script
    [ "$TYPE" != none -a "$DEVICE" != none ] || exit 0
    [ "$DEVICE" != none ] || exit 0
    timeout 60 sh -c "while ! mountpoint -q $MOUNTPOINT && sleep 0.5; do :; done"
end script

post-stop script
    [ -z "$UPSTART_STOP_EVENTS" ] || exit 0
    initctl emit -n mounted MOUNTPOINT="$MOUNTPOINT" TYPE=$TYPE DEVICE=$DEVICE
end script
