# content/etc/init/progress-reboot.conf
#!upstart
description "progress-filesystems-mounted"

start on runlevel 6

task

script

	[ -e /run/nosplash ] && exit 0

	splash --msgtxt="rebooting..." --infinitebar || true 

end script
