# content/etc/init/progress-shutdown.conf
#!upstart
description "progress-shutdown/reboot"

start on starting rc RUNLEVEL=[06]

env RUNLEVEL=0
env MSG="shutting down..."

pre-start script
    start wait-for-state WAIT_FOR=xbmc WAITER=progress-shutdown WAIT_STATE=stopped ACT=no || :

    if [ -e /run/nosplash ]; then
        splash --exit || :
        start xbian-chvt TTYNR=7 || :
    fi

    [ $RUNLEVEL -eq 0 ] || MSG="rebooting..."

    /etc/xbian-initramfs/initram.switcher.sh || :

    [ ! -e /etc/default/xbmc ] || . /etc/default/xbmc
    [ ! -x "$SCROFFRUN" ] || [ $RUNLEVEL -eq 6 ] || MSG='running custom shutdown script...'

    splash --msgtxt="$MSG" --infinitebar --flip || :

    [ ! -x "$SCROFFRUN" ] || [ $RUNLEVEL -eq 6 ] || $SCROFFRUN || :
    [ "$SCROFFHALT" = no ] || echo "standby 0" | cec-client -s -d 1|| :

    if pgrep lightdm >/dev/null; then
        { service lightdm stop; while pgrep lightdm; do sleep 1 && pkill lxsession || sleep 1; done; sleep 5; } || :
    fi
end script

script
    set +e
    [ $RUNLEVEL -eq 0 ] || MSG="rebooting..."

    echo $$ > /run/sendsigs.omit.d/progress-shutdown.pid
    splash --msgtxt="$MSG" --infinitebar --flip || :; ln -s /run/splash/splash.pid /run/sendsigs.omit.d/splash.pid

    exec sleep 365d
end script

pre-stop script
    kill $(cat /run/splash/splash.pid) || :
end script 
