# content/etc/init/progress-reboot.conf
#!upstart
description "progress-reboot"

start on runlevel [6]

task

pre-start exec start xbian-waiter WAIT_FOR=xbmc WAITER=progress-reboot STATUS="stop/waiting"

script
    if [ -e /run/nosplash ]; then
        splash --exit
    else
        splash --msgtxt="rebooting..." --infinitebar --flip || true
    fi

end script
