# content/etc/init/reboot-force.conf
#!upstart

description "should the reboot be stuck from what ever reason, after 100s from reboot/shutdown started do a sync and reboot -f"

start on runlevel [06]

task

pre-start script
    if sleep 60; then sync; reboot -f; fi &  \
        echo $! > /run/sendsigs.omit.d/reboot-force.pid || :
end script
