# content/etc/init/rpcbind.conf
# dummy for sysv services

start on started portmap
stop on stopped portmap
