# content/etc/init/scsi-min-power.conf
start on runlevel [!06]

task

script
    for f in $(find -L /sys/class/scsi_host -maxdepth 2 -name link_power_management_policy 2>/dev/null); do
        echo min_power > $f
    done
end script

