# content/etc/init/setupcon.conf
start on starting xbian-run or stopped xbmc-kill-splash

task

exec setupcon -f
