# content/etc/init/setupcon.conf
start on started xbmc-kill-splash

task

exec setupcon -f --force 
