# content/etc/init/setupcon.conf
start on started xbian-run or stopped xbmc-kill-splash

task

exec setupcon
