# content/etc/init/setupcon.conf
start on started tty1 and stopped rc RUNLEVEL=[!06S] or stopped xbmc-kill-splash

task

exec setupcon -f
