# content/etc/init/setupcon.conf
start on stopped xbmc-kill-splash

task

exec setupcon -f --force 
