# content/etc/init/setupcon.conf
start on runlevel [!06]

task

exec setupcon -f
