# content/etc/init/setupcon.conf
start on started tty1 and stopped rc RUNLEVEL=[!06S]

task

exec start -n xbmc-kill-splash CONSOLE=y
