# content/etc/init/sysklogd.conf
# upstart dummy syslog

start on startup
