# content/etc/init/xbian-run.conf
start on mountall or startup

nice -3
env RAMTMP=

emits mounted

pre-start script
    [ ! -e /etc/default/rcS ] || . /etc/default/rcS ||:
    [ ! -e /etc/default/tmpfs ] || . /etc/default/tmpfs ||:

    __mount() {
        mountpoint -q $1 || \
        case $1 in
            /run)
                [ -n "$RUN_SIZE" ] || RUN_SIZE=10%
                /bin/mount -n -t tmpfs -o size=$RUN_SIZE,nosuid,nodev,noexec,noatime tmpfs /run || :
                ;;
            /tmp)
                [ x"$RAMTMP" = xyes ] || return 0
                [ -z "$TMP_SIZE" ] || TMPFS_SIZE=$TMP_SIZE
                [ -n "$TMPFS_SIZE" ] || TMPFS_SIZE=20%
                /bin/mount -n -t tmpfs -o size="$TMPFS_SIZE",noatime tmpfs /tmp || :
                ;;
        esac
        initctl emit -n mounted MOUNTPOINT=$1
    }

    ### rootfs remount to rw, mount modules xbian
    for f in / /lib/modules; do
        /bin/mount  -n -t xbian $f || :
    done

    [ -d /selinux ] || mkdir -p /selinux
    [ ! -h /etc/mtab ] && { rm -f /etc/mtab; ln -s /proc/mounts /etc/mtab; } || :

    ### mount tmpfs
    for f in /run /tmp; do
        __mount $f
    done

    if ! grep -wq splash /proc/cmdline; then
        touch /run/nosplash
        dmesg -D
        /usr/local/bin/busybox setconsole /dev/tty7 || :
        start xbian-chvt TTYNR=7 || :
    fi
end script

post-start script
    grep -v ^# /etc/fstab| awk '{print $3}' | grep -q 'nfs\|nfs4\|cifs' || initctl emit -n xbian-nowait || :

    [ -d /run/sendsigs.omit.d ] || mkdir -p /run/sendsigs.omit.d
    [ -d /run/lock ] || mkdir -p /run/lock
    [ -d /run/shm ] || mkdir -p /run/shm
    chmod 1777 /run/lock /run/shm /tmp || :
    [ -e /dev/shm ] || ln -s /run/shm /dev/shm
end script
