# content/etc/init/xbian-firewall.conf
start on net-device-up IFACE!=lo or net-device-down IFACE!=lo or started network-interface INTERFACE!=lo or stopped network-interface INTERFACE!=lo

env ENABLEFW=no
env RULEFILE=/run/network/fw.script

pre-start script
    rm -f $RULEFILE
    [ ! -e /etc/default/xbian-firewall ] || . /etc/default/xbian-firewall
    [ "x$ENABLEFW" = xyes -a -x /sbin/iptables ] || { stop; exit 0; }

    echo "# Generated by xbian-firewall
*filter
:INPUT DROP [0:0]
:FORWARD ACCEPT [0:0]
:OUTPUT ACCEPT [1207:92253]
-A INPUT -m state --state RELATED,ESTABLISHED -j ACCEPT" > $RULEFILE

    for ifa in $(ls /sys/class/net); do
        for addr in $(ip addr show $ifa | grep -w 'inet\|inet6' | awk '{print $2}'); do
            case $ifa in
                lo)
                    fwl="-A INPUT -d $addr -i $ifa -j ACCEPT"
                    ;;
                *)
                    fwl="-A INPUT -s $addr -i $ifa -m state --state NEW -j ACCEPT"
                    #### -d ${addr%%/*}
                    ;;
            esac
            echo "$fwl" >> $RULEFILE
        done
    done

    [ -e $RULEFILE ] || exit 0

    echo COMMIT >> $RULEFILE
    cat $RULEFILE | iptables-restore
end script

post-stop script
    for ch in INPUT OUTPUT FORWARD; do iptables -P $ch ACCEPT; done
    iptables -F
end script
