# content/var/tmp/xbian_update/files_to_overwrite/etc/cron.daily/apt-show-versions
#!/bin/sh

exit 0

[ -f /usr/bin/apt-show-versions ] || exit 0

apt-show-versions -i
